# systemd/nstatus-conky.service
[Unit]
Description=NStatus Conky Widget
Documentation=https://github.com/t2o0n321/nstatus-for-ubuntu-desktop
After=nstatus.service graphical-session.target
Requires=nstatus.service
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/conky --config=%h/.config/nstatus/conky/nstatus.conf

# Give the daemon time to write the first state file before Conky starts.
ExecStartPre=/bin/sleep 3

Restart=on-failure
RestartSec=5

# Do NOT hard-code DISPLAY or WAYLAND_DISPLAY here.
# systemd --user inherits the correct display server variables from the
# login session automatically (via PAM / dbus-update-activation-environment).

StandardOutput=journal
StandardError=journal
SyslogIdentifier=nstatus-conky

[Install]
WantedBy=graphical-session.target
